# SOURCES/usb-can-up@.service
[Unit]
Description=VSCOM CAN up %i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/sbin/ip link set up %i
Restart=no
